# ha_admin_mcp/rootfs/etc/services.d/ha-admin-mcp/run
#!/bin/sh
set -eu

echo "[ha-admin-mcp] Starting EXTREMELY DANGEROUS privileged MCP server"
exec python3 /app/src/server.py
